# deb-root/etc/init/graphite-clickhouse.conf
# graphite-clickhouse - Graphite frontend for ClickHouse

description     "Graphite frontend for ClickHouse"

start on runlevel [2345]
stop on runlevel [016]

respawn
respawn limit 3 12

script
	exec start-stop-daemon --start --exec /opt/graphite-clickhouse/graphite-clickhouse
end script
